# agent-rs/packaging/systemd/doro-agent.service
[Unit]
Description=DoroheDoro Rust Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=doro-agent
Group=doro-agent
EnvironmentFile=-/etc/doro-agent/agent.env
ExecStartPre=/usr/bin/doro-agent check-config --config /etc/doro-agent/config.yaml
ExecStart=/usr/bin/doro-agent run --config /etc/doro-agent/config.yaml
Restart=always
RestartSec=5
WorkingDirectory=/var/lib/doro-agent
StateDirectory=doro-agent
RuntimeDirectory=doro-agent
LogsDirectory=doro-agent
UMask=0077
NoNewPrivileges=true
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
